# supervisord.conf
[supervisord]
nodaemon=true

[program:proxy1_token]
command=python3 /app/proxy1/get_token.py --portal %(ENV_PORTAL)s --ch %(ENV_CH_ID)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:proxy1_update]
command=python3 /app/proxy1/update_playlist.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:proxy2]
command=python3 /app/proxy2/main.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:http80]
command=python3 -m http.server 80 --bind 0.0.0.0
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
